# system/hello.service
[Unit]
Description="Hello"

[Service]
ExecStart=/bin/echo "Hello"
Type=simple

[Install]
WantedBy=multi-user.target
